# tetra/templates/systemd/tetra.service
[Unit]
Description=Tetra Configuration at Distance Server
Documentation=https://github.com/your-org/tetra
After=network.target
Wants=network.target

[Service]
Type=simple
User=tetra
Group=tetra
WorkingDirectory=/opt/tetra
Environment=NODE_ENV=production
Environment=TETRA_ENV=production
Environment=PORT=4444
Environment=TETRA_PORT=4444
Environment=TETRA_SRC=/opt/tetra
Environment=TETRA_DIR=/var/lib/tetra
EnvironmentFile=-/etc/tetra/tetra.env
ExecStartPre=/bin/mkdir -p /var/lib/tetra/logs
ExecStartPre=/bin/chown tetra:tetra /var/lib/tetra/logs
ExecStart=/usr/bin/node /opt/tetra/server/server.js
ExecReload=/bin/kill -HUP $MAINPID

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tetra
ReadOnlyPaths=/opt/tetra

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Restart policy
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60s
TimeoutStopSec=30s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tetra

[Install]
WantedBy=multi-user.target