# tetra/templates/systemd/tetra@dev.service
[Unit]
Description=Tetra Service Manager Daemon (@dev environment)
Documentation=https://github.com/pixeljam/tetra
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/tetra

# Environment variables
Environment=TETRA_SRC=/root/src/devops/tetra
Environment=TETRA_DIR=/root/tetra
Environment=NODE_ENV=development
Environment=TETRA_ENV=dev
Environment=PORT=4444
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin

# Load tetra environment and start the tetra service
# This runs: source ~/tetra/tetra.sh && tsm start tetra
ExecStartPre=/bin/bash -c '/bin/mkdir -p /root/tetra/logs /root/tetra/tsm/runtime/processes'
ExecStartPre=/bin/bash -c '/bin/chmod 755 /root/tetra /root/tetra/logs'
ExecStart=/bin/bash -c 'source /root/tetra/tetra.sh && tsm start tetra'
ExecStop=/bin/bash -c 'source /root/tetra/tetra.sh && tsm stop tetra'
ExecReload=/bin/bash -c 'source /root/tetra/tetra.sh && tsm restart tetra'

# Restart policy
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s

# Security settings (relaxed for development)
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tetra-dev

[Install]
WantedBy=multi-user.target
